# docker/etc/supervisord.d/sonarqube.ini
[program:sonarqube]
command=/opt/sonar/bin/linux-%(ENV_ARCH)s/sonar.sh start
stderr_logfile=%(ENV_SONAR_LOG_PATH)s/supervisor.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=1MB
startsecs=15
stopsignal=TERM
stopwaitsecs=30
autorestart=true
startretries=3
